# rootfs/etc/services.d/rspamd/run
#!/command/with-contenv /bin/sh
/usr/bin/rspamd -i -f -u rspamd -g rspamd